# crontab
* * * * *  echo "Cron is OK!" >> /var/log/cron.log 2>&1
#Вот такая джигуржа - нжна пустая строчка

